# infra/systemd/holes-p0-pusher.timer
[Unit]
Description=Push P0 holes every 5 min

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target