# landtek_telegram/landtek-tg-inbox.service
[Unit]
Description=LandTek Telegram inbox receiver (bulletproof webhook)
After=network.target docker.service postgresql.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/landtek
EnvironmentFile=/root/landtek/.env
Environment=PYTHONUNBUFFERED=1
Environment=LANDTEK_TG_PORT=8766
Environment=LANDTEK_TG_PG_DSN=postgresql://n8n:n8npassword@172.18.0.3:5432/n8n
ExecStart=/usr/bin/python3 /root/landtek/landtek_telegram/inbox.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
