# systemd/landtek-scan-intake.service
[Unit]
Description=LandTek ScannerPro → corpus intake (land unclassified → OCR ladder → reenrich → embed)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/landtek
EnvironmentFile=/root/landtek/.env
ExecStart=/usr/bin/python3 /root/landtek/scripts/scan_intake.py --limit 500
# degrade-don't-crash: a bad run must not wedge the timer
SuccessExitStatus=0 1
