# systemd/landtek-proactive.service
[Unit]
Description=LandTek proactive deadline surfacing
After=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/landtek/proactive_deadlines.py
StandardOutput=append:/var/log/landtek_proactive.log
StandardError=append:/var/log/landtek_proactive.log
