# systemd/landtek-verify.service
[Unit]
Description=LandTek verification ladder auto-promote
After=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/landtek/verification_ladder.py
StandardOutput=append:/var/log/landtek_verify.log
StandardError=append:/var/log/landtek_verify.log
