# systemd/landtek-reenrich.service
[Unit]
Description=LandTek post-OCR re-enrichment (SIGNALS v4; O(changed-docs) via text_hash staleness)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/landtek
ExecStart=/usr/bin/python3 /root/landtek/scripts/reenrich.py --limit 2000 --quiet
TimeoutStartSec=1800
