# systemd/landtek-verify.service
[Unit]
Description=LandTek verifiable-corpus loop (scour + point + measure; strict, gate-aligned)
After=network-online.target

# NOTE: this unit formerly ran verification_ladder.py, which auto-promotes entities to
# 'verified' by CORROBORATION (3+ mentions) — that contradicts the provenance write-gate
# (deploy_504: verified = resolving source doc + quoted excerpt). The ladder is quarantined
# (its timer stays disabled) until its target tier is relabeled 'corroborated' so that
# 'verified' uniformly means document-proven. The daily verification job is now verify_loop.py.

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/bin/python3 /root/landtek/scripts/verify_loop.py; /usr/bin/python3 /root/landtek/scripts/doc_discovery.py; /usr/bin/python3 /root/landtek/scripts/contradiction.py; /usr/bin/python3 /root/landtek/scripts/case_dossier.py --all'
StandardOutput=append:/var/log/landtek_verify.log
StandardError=append:/var/log/landtek_verify.log
