# systemd/landtek-digest.timer
[Unit]
Description=Daily LandTek digest at 9 AM Manila (1 AM UTC)

[Timer]
OnCalendar=*-*-* 01:00:00 UTC
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
