# systemd/landtek-expense-harvester.timer
[Unit]
Description=Run the MWK expense harvester daily (keep the reimbursement ledger live)

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
