# systemd/landtek-verify.timer
[Unit]
Description=Verification ladder promotion - daily 3 AM UTC
[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true
[Install]
WantedBy=timers.target
