# systemd/landtek-comm-agent-soak.service
[Unit]
Description=LandTek COMM-AGENT-MAX shadow soak — drive new inbound through L4 (shadow, sends nothing) + monitor per-role hair-split decisions. Read-only w.r.t. production.
After=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=-/root/landtek/.env
WorkingDirectory=/root/landtek
ExecStart=/usr/bin/python3 /root/landtek/scripts/comm_agent_soak.py --tick
StandardOutput=append:/var/log/landtek_comm_agent_soak.log
StandardError=append:/var/log/landtek_comm_agent_soak.log
